# supervisord.conf
[supervisord]
nodaemon=true

[program:emeraldd]
command=/home/emerald/emeraldd
